# Raspberry Pi Program/Gradlog/99-usbhook.rules
ACTION=="add",KERNEL=="sd*", SUBSYSTEMS=="usb", ATTRS{product}=="*", RUN+="/mnt/raid/Gradlog/usbhook.sh %k"
